# recipes-python_st/imx6sx-prp/startprp/prp.service
[Unit]
Description=service start prp
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/bin/sh /PRP.sh
StandardOutput=tty

[Install]
WantedBy=sysinit.target